# deploy/systemd/pet-dashboard-sync.service
[Unit]
Description=Run the pet-dashboard Tractive sync once
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/pet-dashboard
ExecStart=/usr/bin/docker compose --profile sync run --rm sync
